# service/card_proxy_service.cron
SHELL=/bin/bash
*/5 * * * * cpr_service card_proxy_cleaner.sh
