# service/card_proxy_service.cron
SHELL=/bin/sh
* * * * * cpr_service card_proxy_keyproc.sh
* * * * * cpr_service card_proxy_hmacproc.sh
* * * * * cpr_service card_proxy_cleaner.sh
