# debian/card-proxy-confpatch.postinst
#!/bin/bash

CP_USER="cpr_confpatch"

case "$1" in
    "configure")
        grep -q "^$CP_USER:" /etc/passwd || ( \
            adduser --gecos 'CardProxy.ConfPatch' --disabled-password $CP_USER ; \
            mkdir -p /var/run/$CP_USER /var/log/$CP_USER ; \
            chown $CP_USER:root /var/log/$CP_USER /var/run/$CP_USER )

        ln -s ../sites-available/10-card_proxy_confpatch /etc/nginx/sites-enabled
        service nginx restart
        service card_proxy_confpatch restart

        update-rc.d card_proxy_confpatch defaults
        ;;

    "upgrade" | "abort-upgrade" | "abort-remove" | "abort-deconfigure")

        service card_proxy_confpatch restart

        ;;

    *) 
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0

